# packaging/DEBIAN/prerm
#!/bin/bash
set -e

# Stop and disable service before removal
if systemctl is-active --quiet hyperweaver-agent; then
    echo "Stopping Hyperweaver Agent service..."
    systemctl stop hyperweaver-agent
fi

if systemctl is-enabled --quiet hyperweaver-agent; then
    echo "Disabling Hyperweaver Agent service..."
    systemctl disable hyperweaver-agent
fi

#DEBHELPER#
